# picam-stream.service
[Unit]
Description=PiCam Stream — Live Camera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rp
Group=rp
WorkingDirectory=/home/rp/camera_web-test-1
ExecStart=/bin/bash /home/rp/camera_web-test-1/start.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Give time for git pull + pip install
TimeoutStartSec=120

# Environment
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
